# crontab
# Hourly db backup
36 * * * * YOUR_MISSKEY_KIT_DIR/backup-db.sh

